# packaging/systemd/wesnothd.service.in
[Unit]
Description=Wesnoth Multiplayer Server Daemon
After=network.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/wesnothd

[Install]
WantedBy=multi-user.target
